# aws_monitoring/crontab
0 10 * * * cd /app && python manage.py check >> /var/log/cron.log 2>&1
