# aws_monitoring/crontab
0 10 * * * root cd /app && /usr/local/bin/python manage.py check >> /var/log/cron.log 2>&1
